# files/etc/systemd/system/node_exporter.socket.jinja2
[Unit]
Description=Node Exporter

[Socket]
ListenStream={{ http.get('addr') }}:{{ http.get('port') }}

[Install]
WantedBy=sockets.target
